# docker/container/supervisord.conf
# Process model for the all-in-one Portal image. supervisord (PID-managed under tini) runs the
# five long-lived services + the internal Caddy router. Each program logs to the container's
# stdout (supervisord captures child output and writes it as root, so the child's own user
# doesn't matter) and restarts on crash. Startup ordering: Postgres + Redis come up first
# (priority 10); the api/worker wrappers wait for them before serving. Schema migrations + the
# admin seed already ran in the entrypoint, so the app processes never race an unmigrated DB.
#
# HOME=/app on the DB-connecting programs: asyncpg reads a passfile under $HOME while parsing
# connect args. Without this, a child running as uid `app` would inherit supervisord's HOME=/root,
# which `app` can't traverse → EACCES before it ever reaches the socket. /app is world-traversable.
[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/run/supervisord.pid

[unix_http_server]
file=/run/supervisor.sock

[supervisorctl]
serverurl=unix:///run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:postgres]
command=/usr/lib/postgresql/16/bin/postgres -D /data/postgres -c listen_addresses=127.0.0.1
user=postgres
priority=10
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:redis]
command=redis-server --appendonly yes --dir /data/redis --bind 127.0.0.1
user=redis
priority=10
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:api]
command=/usr/local/bin/portal-start-api.sh
directory=/app
user=app
environment=HOME="/app"
priority=20
startsecs=10
startretries=5
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:worker]
# Runs as root because it writes synced files to the host mount, which may squash non-root uids
# (e.g. an NFS export). It has no inbound network, so this is low risk; the internet-facing web
# and api stay non-root.
command=/usr/local/bin/portal-start-worker.sh
directory=/app
user=root
environment=HOME="/app"
priority=20
startsecs=10
startretries=5
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:web]
command=node /app/web/build
directory=/app/web
user=app
priority=20
environment=PORT="3000",HOST="127.0.0.1",HOME="/app"
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:caddy]
command=caddy run --config /etc/caddy/Caddyfile --adapter caddyfile
user=root
priority=30
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
